# deploy/geo-celery.service
[Unit]
Description=GEO Explorer Celery Worker
After=network.target geo-redis.service
Requires=geo-redis.service

[Service]
Type=simple
User=ffh
WorkingDirectory=/home/ffh/geo-explorer
Environment=PATH=/home/ffh/geo-explorer/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/ffh/geo-explorer/.venv/bin/python -m celery -A src.celery_app worker --pool=solo --loglevel=info
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
